# modules/obscuring/airootfs.any/etc/systemd/system/onion-grater.service
[Unit]
Description=Tor control port filter proxy
Documentation=https://tails.boum.org/contribute/design/

[Service]
Type=simple
ExecStart=python /usr/lib/onion-grater

[Install]
WantedBy=multi-user.target

